# mocks/services/lipa-api/run
#!/bin/sh
# Generated by Chef
# Twiket Extranet lipa-api service

exec 2>&1
[ -d env ] && setenv="-e $(pwd)/env"

exec chpst -o 65000 -P -u twiket $setenv /home/twiket/lipa/current/api -c /etc/twiket/lipa/lipa-api.yml run
