# mocks/services/lipa-workers/run
#!/bin/sh
# Generated by Chef
# Twiket lipa-workers service

exec 2>&1
[ -d env ] && setenv="-e $(pwd)/env"

exec chpst -o 65000 -P -u twiket $setenv /home/twiket/lipa/current/workers -c /etc/twiket/lipa/lipa-workers.yml run
