# mocks/services/twapi-lira/run
#!/bin/sh
# Generated by Chef

exec 2>&1

[ -d env ] && setenv="-e $(pwd)/env"

export NODE_CONFIG_DIR=./mocks/services/lira/current/config

if [ -f n_version ]; then
  node_version=$(cat n_version)
  n $node_version
  options=""
  if [ -f n_options ]; then
    options=$(cat n_options)
    echo "{\"@timestamp\": \"$(date +"%FT%T.%3NZ")\", \"message\": \"starting lira with options: $options\"}"
  else
    echo "{\"@timestamp\": \"$(date +"%FT%T.%3NZ")\", \"message\": \"n_options not found, continue...\"}"
  fi
  exec chpst -u twiket $setenv n use $node_version $options index.js
else
  echo "{\"@timestamp\": \"$(date +"%FT%T.%3NZ")\", \"message\": \"n_version not found for lira\"}"
fi