# mocks/services/pipa/run
#!/bin/sh
# Generated by Chef
# Twiket meta service


exec 2>&1
[ -d env ] && setenv="-e $(pwd)/env"

sv start twapi-refs
while true; do
  if nc -z 127.0.0.1 5050; then
    echo 'twapi-refs runned'
    break
  else
    echo 'twapi-refs still runnning'
    sleep 1
  fi
done

sv start twapi-cache
while true; do
  if nc -z 127.0.0.1 5060; then
    echo 'twapi-cache runned'
    break
  else
    echo 'twapi-cache still runnning'
    sleep 1
  fi
done

sv start twapi-main-3031

#Webfront Notification

for i in $(seq 0 100); do
  if curl; then
    break
  else
    echo 'pipa is not running. Sleeping 3 sec'
    sleep 3
  fi
done

exec sleep infinity